# src/etc/cron.d/00-oc-logrotate
#Ensure logrotate is run every 5 mins.
*/5 * * * * root /etc/cron.daily/logrotate >> /dev/null
